# .infrastructure/docker-compose/crontab
SHELL=/bin/bash

12 * * * * docker exec php-fpm bin/console app:currency:updateRates >/dev/null 2>&1